# hub-server/src/main/dist/hub-server.service
[Unit]
Description=Datawire Hub Server
After=network.target
ConditionPathExists=/opt/datawire/hub-server

[Service]
ExecStart=/opt/datawire/hub-server/bin/hub-server /opt/datawire/hub-server/prod.yml
Type=forking
Restart=on-failure

[Install]
WantedBy=multi-user.target