# dist/roles/gateway/files/discovery-gateway.service
[Unit]
Description=Datawire Discovery Gateway
After=cloud-final.service
ConditionPathExists=/opt/datawire/discovery-gateway

[Service]
ExecStart=/opt/datawire/discovery-gateway/bin/discovery-gateway server /opt/datawire/discovery-gateway/config.yml
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target