# discovery-server/src/main/dist/discovery-server.service
[Unit]
Description=Datawire Discovery Server
After=network.target
ConditionPathExists=/opt/datawire/discovery-server

[Service]
ExecStart=/opt/datawire/discovery-server/bin/discovery-server server /opt/datawire/discovery-server/prod.yml
Type=forking
Restart=on-failure

[Install]
WantedBy=multi-user.target