# discovery-web/dist/systemd-discovery.service
# file: discovery.conf

[Unit]
Description=Datawire Discovery
After=network.target
AssertPathExists=/opt/discovery
AssertFileIsExecutable=/opt/discovery/bin/discovery.sh

[Service]
Type=simple
EnvironmentFile=/etc/discovery/discovery.env
ExecStart=/opt/discovery/bin/discovery.sh

[Install]
WantedBy=multi-user.target