# discovery-web/dist/systemd-discovery.service
# file: discovery.conf

[Unit]
Description=Datawire Discovery
After=cloud-final.service enventurer.service
AssertPathExists=/opt/discovery
AssertFileIsExecutable=/opt/discovery/discovery.sh

[Service]
Type=simple
EnvironmentFile=/etc/datawire/environment
ExecStart=/opt/discovery/discovery.sh
WorkingDirectory=/opt/discovery

[Install]
WantedBy=multi-user.target