# dist/roles/discovery/files/discovery-server.service
[Unit]
Description=Datawire Discovery Server
After=cloud-final.service
ConditionPathExists=/opt/datawire/discovery-server

[Service]
ExecStart=/opt/datawire/discovery-server/bin/discovery-server server /opt/datawire/discovery-server/config.yml
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target