# hub-gateway/src/dist/hub-gateway.service
[Unit]
Description=Datawire Hub Gateway
After=network.target
ConditionPathExists=/opt/datawire/hub-gateway

[Service]
ExecStart=/opt/datawire/hub-gateway/bin/hub-gateway server /opt/datawire/hub-gateway/prod.yml
Type=forking
Restart=on-failure

[Install]
WantedBy=multi-user.target